# backend/systemd/tax-filing-celery.service
[Unit]
Description=Vietnamese Tax Filing Celery Worker
After=network.target redis.service
Wants=redis.service

[Service]
Type=exec
User=taxapp
Group=taxapp
WorkingDirectory=/home/taxapp/tax-filing-app
Environment=PATH=/home/taxapp/tax-filing-app/venv/bin
ExecStart=/home/taxapp/tax-filing-app/venv/bin/python run_celery.py
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tax-filing-celery

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/home/taxapp/tax-filing-app/uploads /home/taxapp/tax-filing-app/logs

[Install]
WantedBy=multi-user.target
